# backend/crontab
#Collect draft data every two days
0 */12 * * * python CollectDraftData.py
